# supervisor/conf.d/gunicorn.conf
[program:gunicorn]
directory=/opt/infomarquee
command=/usr/bin/gunicorn -b 0.0.0.0:8000 app:app
